# deploy/podman/flightctl-pam-issuer/flightctl-pam-issuer.container
[Unit]
Description=Flight Control PAM OIDC Issuer server
PartOf=flightctl.target
After=flightctl-api-certs-init.service
Wants=flightctl-api-certs-init.service

[Container]
ContainerName=flightctl-pam-issuer
Image={{ .PamIssuer.Image }}:{{ .PamIssuer.Tag }}
Pull=newer
Network=flightctl.network

PublishPort=8444:8444
Volume=/etc/flightctl/pki:/root/.flightctl/certs:ro,z
Volume=flightctl-pam-issuer-etc:/etc:rw
Volume=/etc/flightctl/flightctl-pam-issuer/config.yaml:/root/.flightctl/config.yaml:ro,z

[Service]
Restart=always
RestartSec=30

ExecStartPre=/usr/bin/flightctl-standalone render template --input-file /usr/share/flightctl/flightctl-pam-issuer/config.yaml.template --output-file /etc/flightctl/flightctl-pam-issuer/config.yaml

[Install]
WantedBy=flightctl.target
