# deploy/podman/flightctl-db/flightctl-db-standalone.container
[Unit]
Description=Flight Control Database service

[Container]
ContainerName=flightctl-db
Environment=POSTGRESQL_DATABASE=flightctl POSTGRESQL_USER=flightctl_app POSTGRESQL_MASTER_USER=admin
Image=quay.io/sclorg/postgresql-16-c9s:20250214
Pull=missing
PublishPort=5432:5432
Volume=flightctl-db:/var/lib/pgsql/data:Z
Secret=flightctl-postgresql-password,type=env,target=PGPASSWORD
Secret=flightctl-postgresql-master-password,type=env,target=POSTGRESQL_MASTER_PASSWORD
Secret=flightctl-postgresql-user-password,type=env,target=POSTGRESQL_PASSWORD
# Note: Migrator password removed to prevent credential leakage - migrator user creation
# is now handled by the migration service which has its own credentials
Volume=/usr/share/flightctl/flightctl-db/enable-superuser.sh:/usr/share/container-scripts/postgresql/start/enable-superuser.sh

[Service]
Restart=always
RestartSec=30
