# deploy/podman/flightctl-worker/flightctl-worker.container
[Unit]
Description=FlightCtl Worker service
After=flightctl-db.service
Requires=flightctl-db.service

[Container]
ContainerName=flightctl-worker
Image=quay.io/flightctl/flightctl-worker:0.4.0
Network=flightctl.network
Environment=HOME=/root
Volume=%h/.config/flightctl/flightctl-worker-config/config.yaml:/root/.flightctl/config.yaml:Z

[Service]
Restart=always
RestartSec=30
Slice=flightctl.slice

[Install]
WantedBy=flightctl.slice
