# deploy/podman/flightctl-imagebuilder-worker/flightctl-imagebuilder-worker.container
[Unit]
Description=Flight Control ImageBuilder Worker service
PartOf=flightctl.target
After=flightctl-db-wait.service flightctl-db-migrate.service flightctl-kv.service flightctl-certs-init.service
Wants=flightctl-db-wait.service flightctl-db-migrate.service flightctl-kv.service flightctl-certs-init.service

[Container]
ContainerName=flightctl-imagebuilder-worker
Image={{ .ImagebuilderWorker.Image }}:{{ .ImagebuilderWorker.Tag }}
Pull=missing
Network=flightctl.network
AddHost=%H:host-gateway
Environment=HOME=/root
Secret=flightctl-postgresql-user-password,type=env,target=DB_PASSWORD
Secret=flightctl-kv-password,type=env,target=KV_PASSWORD
Environment=DB_USER=flightctl_app
Volume=/etc/flightctl/flightctl-imagebuilder-worker/config.yaml:/root/.flightctl/config.yaml:ro,z
Volume=/etc/flightctl/pki/flightctl-api/client-signer.crt:/root/.flightctl/certs/client-signer.crt:ro,z
Volume=/etc/flightctl/pki/flightctl-api/client-signer.key:/root/.flightctl/certs/client-signer.key:ro,z
Volume=/etc/flightctl/pki/ca-bundle.crt:/root/.flightctl/certs/ca-bundle.crt:ro,z
Volume=/etc/flightctl/pki/db:/root/.flightctl/certs/db:ro,z
Volume=/var/tmp/flightctl-builds:/var/tmp/flightctl-builds:rw,z
Volume=/var/tmp/flightctl-exports:/var/tmp/flightctl-exports:rw,z
Volume=flightctl-worker-storage:/var/lib/containers
# Optional: Mount RHEL entitlement certs for subscription-only repos (auto-detected if present).
# On RHEL hosts with active subscription, uncomment to use host entitlements:
# Volume=/etc/pki/entitlement:/etc/pki/entitlement:ro,z
PodmanArgs=--privileged

[Service]
Restart=always
RestartSec=30
Delegate=yes

ExecStartPre=/usr/bin/flightctl-standalone render template --input-file /usr/share/flightctl/flightctl-imagebuilder-worker/config.yaml.template --output-file /etc/flightctl/flightctl-imagebuilder-worker/config.yaml

[Install]
WantedBy=flightctl.target
