# deploy/podman/flightctl-db-migrate/flightctl-db-wait.container
[Unit]
Description=Flight Control Database Wait (admin)
PartOf=flightctl.target
After=flightctl-db.service 

[Container]
ContainerName=flightctl-db-wait
Image=quay.io/flightctl/flightctl-db-setup:latest
Network=flightctl.network
AutoUpdate=none
Environment=DB_HOST=flightctl-db
Environment=DB_PORT=5432
Environment=DB_NAME=flightctl
Environment=DB_USER=admin
Secret=flightctl-postgresql-master-password,type=env,target=DB_PASSWORD
Exec=/app/deploy/scripts/wait-for-database.sh --timeout=60 --sleep=1

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=flightctl.target


