# deploy/podman/flightctl-db-migrate/flightctl-db-users-init.container
[Unit]
Description=Flight Control Database Users Initialization
PartOf=flightctl.target
After=flightctl-db-wait.service
Wants=flightctl-db-wait.service

[Container]
ContainerName=flightctl-db-users-init
Image={{ .DbSetup.Image }}:{{ .DbSetup.Tag }}
Pull=missing
Network=flightctl.network
AutoUpdate=none
Environment=DB_HOST=flightctl-db
Environment=DB_PORT=5432
Environment=DB_NAME=flightctl
Environment=DB_APP_USER=flightctl_app
Environment=DB_MIGRATION_USER=flightctl_migrator
Environment=DB_ADMIN_USER=admin

Secret=flightctl-postgresql-master-password,type=mount,target=/run/secrets/db-admin/masterPassword
Secret=flightctl-postgresql-migrator-password,type=mount,target=/run/secrets/db-migration/migrationPassword
Secret=flightctl-postgresql-user-password,type=mount,target=/run/secrets/db/userPassword

Exec=./deploy/scripts/setup_database_users.sh --direct

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
RestartSec=5
ExecCondition=/bin/bash -c 'test "$(python3 /usr/share/flightctl/yaml_helpers.py extract ".db.type" /etc/flightctl/service-config.yaml --default "builtin")" != "external"'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=flightctl.target
