# test/scripts/agent-images/embedded-demo-app.container
[Container]
Image=quay.io/flightctl-tests/alpine:v1
Exec=sleep infinity

[Install]
WantedBy=default.target
