# deploy/podman/flightctl-alert-exporter/flightctl-alert-exporter.container
[Unit]
Description=Flight Control Alert Exporter service
PartOf=flightctl.target
After=flightctl-db-wait.service flightctl-db-migrate.service flightctl-kv.service flightctl-alertmanager.service
Wants=flightctl-db-wait.service flightctl-db-migrate.service flightctl-kv.service flightctl-alertmanager.service

[Container]
ContainerName=flightctl-alert-exporter
Image={{ .AlertExporter.Image }}:{{ .AlertExporter.Tag }}
Pull=missing
Network=flightctl.network
AddHost=%H:host-gateway
Environment=HOME=/root
Secret=flightctl-postgresql-user-password,type=env,target=DB_PASSWORD
Secret=flightctl-kv-password,type=env,target=KV_PASSWORD
Environment=DB_USER=flightctl_app
Volume=/etc/flightctl/flightctl-alert-exporter/config.yaml:/root/.flightctl/config.yaml:ro,z
Volume=/etc/flightctl/pki/db:/root/.flightctl/certs/db:ro,z

[Service]
Restart=always
RestartSec=30

ExecStartPre=/usr/bin/flightctl-standalone render template --input-file /usr/share/flightctl/flightctl-alert-exporter/config.yaml.template --output-file /etc/flightctl/flightctl-alert-exporter/config.yaml

[Install]
WantedBy=flightctl.target
