# packaging/observability/flightctl-userinfo-proxy.container.template
[Unit]
Description=Flight Control UserInfo Proxy for OAuth Integration
After=flightctl-network.service
Wants=flightctl-network.service
PartOf=flightctl-observability.target

[Container]
Image=${USERINFO_PROXY_IMAGE}
Pull=missing

ContainerName=flightctl-userinfo-proxy

# Network configuration
Network=flightctl

# Environment variables for configuration
Environment=USERINFO_LISTEN_PORT=8080
Environment=USERINFO_UPSTREAM_URL=${USERINFO_UPSTREAM_URL}
Environment=USERINFO_SKIP_TLS_VERIFY=${USERINFO_SKIP_TLS_VERIFY}

# Security
User=1001:1001

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target 